# systemd/certbot.service
[Unit]
Description=Certbot automatic renewal

[Service]
ExecStart=/usr/bin/certbot renew --quiet --deploy-hook "systemctl reload nginx"
